# borgmatic/cron.d/borgmatic.j2
# This file is generated by Ansible
#  DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
#
MAILTO=root
0 * * * * {{ borgmatic_user_name }} {{ borgmatic_pip_virtualenv }}/bin/borgmatic
