# nodepool/elements/nodepool-base/sshd-keygen.target
[Unit]
Wants=sshd-keygen@ecdsa.service
PartOf=sshd.service
